# ansible/roles/pangolin_vps/templates/pangolin-docker.service.j2
# Pangolin Docker Stack systemd service
# Generated by Ansible
# Ensures the Docker Compose stack starts on boot and restarts on failure

[Unit]
Description=Pangolin VPN Server (Docker Compose)
Documentation=https://github.com/fosrl/pangolin
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ pangolin_base_dir }}

# Start the stack
ExecStart=/usr/bin/docker compose up -d --remove-orphans

# Stop the stack gracefully
ExecStop=/usr/bin/docker compose down

# Reload configuration (restart containers)
ExecReload=/usr/bin/docker compose restart

# Environment
Environment=COMPOSE_HTTP_TIMEOUT=120
Environment=DOCKER_CLIENT_TIMEOUT=120

# Restart policy
Restart=on-failure
RestartSec=30s

# Resource limits
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
